# skeleton/inittab
::sysinit:/bin/mount -n -t devtmpfs devtmpfs /dev
::sysinit:/bin/mount -n -t proc proc /proc
::sysinit:/bin/mount -n -t sysfs sysfs /sys
::sysinit:/bin/mount -n -t tmpfs tmpfs /run
::sysinit:/bin/rm -rf /var/tmp/*
::sysinit:/bin/hostname -F /etc/hostname
::sysinit:/sbin/netsetup eth0
tty1::respawn:/sbin/getty 38400 /dev/tty1 linux
tty2::respawn:/sbin/getty 38400 /dev/tty2 linux
tty3::respawn:/sbin/getty 38400 /dev/tty3 linux
tty4::respawn:/sbin/getty 38400 /dev/tty4 linux
